# scripts/xarcade2jstick.service
[Unit]
Description=Xarcade2Jstick

[Service]
Type=simple
ExecStart=/usr/local/bin/xarcade2jstick -s

[Install]
WantedBy=multi-user.target
